# environments/app/crontab
# do daily/weekly/monthly maintenance
# min	hour	day	month	weekday	command
*/1	*	*	*	*	cd /var/www/app && php artisan schedule:run >> /dev/null 2>&1

